# carbon-cache/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid

[program:carbon]
command=/usr/bin/start_carbon.sh
auto_start=true
autorestart=true
