# supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
auto_start=true
autorestart=true
